# deploy/systemd/cdsync@.service
[Unit]
Description=CDSync runner for connection %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cdsync
Group=cdsync
WorkingDirectory=/var/lib/cdsync
ExecStart=/usr/local/bin/cdsync run --config /etc/cdsync/config.yaml --connection %i
Restart=always
RestartSec=5
TimeoutStopSec=120
KillSignal=SIGTERM
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
